# systemd/services/usb-gadget.service
[Unit]
Description=USB Ethernet Gadget
After=network.target
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb-gadget-eth.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
